# aula-11-10-15/echo_uds_activated_service/echod.socket
[Unit]
Description=Echo daemon activation socket

[Socket]
ListenStream=/opt/echod/echod.sock

[Install]
WantedBy=sockets.target
