# aula-11-10-15/echo_uds_activated_service/echod_auto.socket
[Unit]
Description=Echo daemon activation socket

[Socket]
ListenStream=/opt/echod_auto/echod.sock

[Install]
WantedBy=sockets.target
